# exes/ws-bridge/systemd/confman.socket
[Unit]
Description=Confman service socket

[Socket]
ListenStream=/run/tfc/confman.rpc.sock
Accept=yes

[Install]
WantedBy=sockets.target
